# etc/supervisor/conf.d/redis.conf
[program:redis]
autorestart = false
autostart = true
command = redis-server
stderr_logfile = /var/log/redis.err.log
stdout_logfile = /var/log/redis.out.log
